# etc/cron.monthly/meet-get-mitca-cert
#!/bin/sh

( cd /etc/ssl/certs && wget -q https://ca.mit.edu/mitca.crt )
